# supervisord.conf
# ============================================================
# DD Planner — Supervisord Config (Production Container)
#
# RULES — DO NOT CHANGE:
#   - backend priority=10, nginx priority=20
#     Backend MUST start before nginx or cold requests get 502
#   - uvicorn must use port 8001 (nginx proxies here)
#   - workers=1 is intentional for Cloud Run reliability
#   - startretries=5 handles transient startup failures
# ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Backend starts FIRST (priority=10 < nginx priority=20)
[program:backend]
# port 8001 — must match nginx proxy_pass; workers 1 is intentional for Cloud Run reliability
command=uvicorn server:app --host 0.0.0.0 --port 8001 --workers 1 --timeout-keep-alive 120
directory=/app/backend
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PLAYWRIGHT_BROWSERS_PATH="/pw-browsers"
priority=10                      # MUST be lower number than nginx (starts first)

# Nginx starts SECOND (after backend is up)
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stopwaitsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20                      # MUST be higher number than backend (starts second)
